# packer/files/devenv-boot.service
[Unit]
Description=DevEnv first-boot provisioning
After=cloud-init.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/devenv-boot.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
